# mqtt_io.conf
[program:mqtt-io]
command = /usr/bin/python3 -m mqtt_io /home/pi/GPIOconfig.yml
directory = /home/pi
redirect_stderr = true
stdout_logfile = /var/log/mqtt-io.log